# ansible/init-distributor-node/templates/livestreaming-distributor-node.service.j2
[Unit]
Description=Livestreaming distributor node
After=network.target

[Service]
WorkingDirectory=/opt/livestreaming-server
ExecStart=/opt/livestreaming-server/venv/bin/python3 -m livestreaming distributor
User=livestreaming
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target
